# .config/autostart/fusuma.desktop
[Desktop Entry]
Name=Fusuma
Type=Application
Exec=/usr/bin/fusuma -d
OnlyShowIn=GNOME
NoDisplay=true

