# .config/autostart/fusuma.desktop
[Desktop Entry]
Exec=/usr/bin/fusuma
Icon=
Name=fusuma
NoDisplay=true
Path=
Terminal=False
Type=Application
