# .config/autostart/goldendict.desktop
[Desktop Entry]
Name=Goldendict
Type=Application
Exec=/usr/bin/goldendict
OnlyShowIn=GNOME
NoDisplay=true

